# supervisord.conf
[supervisord]
nodaemon=true

[program:xray]
command=/usr/local/bin/xray run -config /usr/local/etc/xray/config.json
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:ttyd]
command=/bin/sh -lc 'if [ -n "$TTYD_USER" ] && [ -n "$TTYD_PASS" ]; then exec /usr/bin/ttyd -W -p 7681 -c "$TTYD_USER:$TTYD_PASS" bash; else exec /usr/bin/ttyd -W -p 7681 bash; fi'
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
